# packaging/host/postrm
#!/bin/sh
# Runs after the host package is removed or purged.
set -e

case "$1" in
  remove)
    if [ -d /run/systemd/system ]; then
      systemctl daemon-reload || true
    fi
    ;;

  purge)
    if [ -x /usr/bin/deb-systemd-helper ]; then
      deb-systemd-helper purge easyhomeserver.service >/dev/null || true
      deb-systemd-helper unmask easyhomeserver.service >/dev/null || true
    fi

    # Only on purge. `apt remove` must not destroy the admin password and settings — a plain
    # remove/reinstall cycle should not silently reset the machine to first-run setup.
    rm -rf /var/lib/easyhomeserver

    if [ -d /run/systemd/system ]; then
      systemctl daemon-reload || true
    fi
    ;;

  upgrade | failed-upgrade | abort-install | abort-upgrade | disappear) ;;

  *)
    echo "postrm called with unknown argument '$1'" >&2
    exit 1
    ;;
esac

exit 0
